# shelly-manager/rootfs/etc/services.d/shelly-api/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Shelly Manager API
# ==============================================================================

# Get addon configuration from Home Assistant
SECRET_KEY=$(bashio::config 'secret_key')
TIMEOUT=$(bashio::config 'timeout')
MAX_WORKERS=$(bashio::config 'max_workers')
LOG_LEVEL=$(bashio::config 'log_level')

bashio::log.info "Starting Shelly Manager API..."

export SHELLY_SECRET_KEY="${SECRET_KEY}"
export SHELLY_NETWORK__TIMEOUT="${TIMEOUT}"
export SHELLY_NETWORK__MAX_WORKERS="${MAX_WORKERS}"
export PATH="/app/.venv/bin:$PATH"
export PYTHONDONTWRITEBYTECODE=1
export PYTHONUNBUFFERED=1

cd /app
exec s6-setuidgid app /app/.venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000 --log-level ${LOG_LEVEL}
