# shelly-manager/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Nginx Web Server
# ==============================================================================

bashio::log.info "Starting nginx..."

# Test nginx configuration
bashio::log.info "Testing nginx configuration..."
nginx -t

bashio::log.info "Starting nginx daemon..."
exec nginx -g 'daemon off;'
